# install/systemd.socket
[Socket]
ListenStream=/run/sail.socket
SocketMode=0660
SocketUser=root
SocketGroup=sail

[Install]
WantedBy=sockets.target
