# install/systemd.socket
[Socket]
ListenStream=/run/sailor.socket
SocketMode=0660
SocketUser=root
SocketGroup=sailors

[Install]
WantedBy=sockets.target
